# files/system/wifi-connect.service
[Unit]
Description=Auto Wi-Fi Connect from /etc/wifi.conf
After=network.target NetworkManager.service
Wants=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/wifi_connect.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target