# files/90-canfd.rules
ACTION=="add", KERNEL=="can0", RUN+="/usr/local/bin/canfd-setup.sh"
